# example.upstart.conf
description "Cachet Monitor"

start on startup

env USER=root
env HOME=/root

setuid root
setgid root
chdir /root

script
	exec cachet-monitor -c /cachet-monitor.json --immediate
end script
